# docker-webhook.service
[Service]
ExecStart=/home/ubuntu/Docker-Continuous-Integration/run.sh
Restart=on-failure
WorkingDirectory=/home/ubuntu/Docker-Continuous-Integration